# avro-test/producer/crontab
#run python script every minutes
*/1 * * * * /usr/local/bin/python /usr/src/app/avro-producer.py > /proc/1/fd/1 2>/proc/1/fd/2
